# com.mihai.beeperclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mihai.beeperclaw</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your virtual environment python executable -->
        <string>/Volumes/990Pro2TB/MyProjects/beeperclaw/venv/bin/beeperclaw</string>
        <string>run</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Volumes/990Pro2TB/MyProjects/beeperclaw</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/beeperclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/beeperclaw.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
